# periodic/15min/docker-prune
#!/bin/sh
docker system prune -f --volumes
