# periodic/15min/docker-prune
#!/bin/sh
docker system prune -f --filter "until=24h"
docker volume prune -f
