# services.d/jenkins/run
#!/usr/bin/execlineb -P
with-contenv
s6-setuidgid jenkins
s6-env HOME=/var/jenkins_home
/usr/local/bin/jenkins.sh
